# apps/backend/scripts/makerhub-backend.service
[Unit]
Description=MakerHub FastAPI Backend
After=network.target postgresql.service
Wants=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/srv/makerhub/apps/backend
Environment=UV_CACHE_DIR=/srv/makerhub/.uv-cache
EnvironmentFile=/srv/makerhub/apps/backend/.env
ExecStart=/root/.local/bin/uv run uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5
TimeoutStartSec=30

[Install]
WantedBy=multi-user.target
